# debian/dsme-dbus.init
#!/sbin/openrc-run

description="Connect dsme to dbus"

depend() {
	after dbus dbus-actdead
}

start() {
	/usr/sbin/waitdbus system
	/usr/sbin/dsmetool --start-dbus
}
