# debian/dsme-thermal.service
[Unit]
Description=Select thermal limits

[Service]
ExecStart=/usr/lib/dsme/dsme-thermal
