# debian/dsme-dbus.service
[Unit]
Description=Connect DSME to D-Bus

[Service]
ExecStart=/usr/sbin/dsmetool --start-dbus
