# debian/dsme-dbus.init
#!/sbin/openrc-run
### BEGIN INIT INFO
# Provides:          dsme-dbus
# Required-Start:    $remote_fs dbus
# Required-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:
# Short-Description: Connect DSME to D-Bus
### END INIT INFO

start() {
	/usr/sbin/waitdbus system
	/usr/sbin/dsmetool --start-dbus
}
