# debian/dsme-thermal.init
#!/sbin/openrc-run

description="Select thermal limits"

depend() {
	after dbus dbus-actdead
}

start() {
	/usr/lib/dsme/dsme-thermal
}
